# usr/profile.d/10-aliases.ps1
Set-Alias scoopupdate "$HOME\.local\bin\scoop-interactive-update.ps1"
