# Documents/PowerShell/profile.d/99-command-not-found.ps1
#f45873b3-b655-43a6-b217-97c00aa0db58 PowerToys CommandNotFound module
# PowerShell 7+：使用 PowerToys 模組（需安裝 PowerToys 並啟用 CommandNotFound）
Import-Module -Name Microsoft.WinGet.CommandNotFound -ErrorAction SilentlyContinue
#f45873b3-b655-43a6-b217-97c00aa0db58
